# 09-bash/etc/systemd/system/ws_status.service
[Unit]
Description=ws_status
Requires=nginx.service
Requires=sendmail.service
[Service]
Type=oneshot
ExecStart=/opt/ws_status.sh
